# crontab
* * * * *  /var/www/html/rehan_school/artisan schedule:run >> /dev/null 2>&1

